# roles/systemd_timer/templates/etc/systemd/system/template.timer.j2
[Unit]
Description=Timer: {{ timer.description }} (sheduled as timer)

[Timer]
OnBootSec={{ timer.interval_sec }}
OnUnitActiveSec={{ timer.interval_sec }}
Unit={{ timer.name }}.service

[Install]
WantedBy=timers.target
